# mkosi.extra/usr/lib/systemd/system/mkfs.service
[Service]
Type=oneshot
ExecStart=mkfs.btrfs -f /dev/nvme0n1
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
