# mkosi.profiles/cxl/mkosi.extra/usr/lib/systemd/system/cxl.service
[Service]
Type=oneshot
ExecStart=/bin/bash -c '/usr/sbin/cxl_probe.sh'
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
